# infra/podman/quadlet/indomarket-api.container
[Unit]
Description=IndoMarket Insight API Container
After=network-online.target
Wants=network-online.target

[Container]
Image=localhost/indomarket-insight:latest
ContainerName=indomarket-api
Network=indomarket.network
PublishPort=8000:8000
Volume=%h/indomarket_insight/data:/app/data:Z
Volume=%h/indomarket_insight/models:/app/models:Z
Volume=%h/indomarket_insight/logs:/app/logs:Z
Volume=%h/indomarket_insight/config:/app/config:Z
Volume=%h/indomarket_insight/backups:/app/backups:Z
Environment=TZ=Asia/Jakarta
Exec=uvicorn app.api:api --host 0.0.0.0 --port 8000

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=default.target
